# .crontab
# Min   Hour    mday    month   wday(0=sun)     command

0	*	*	*	*		/usr/sbin/anacron -t "$HOME/.anacron/anacrontab" -S "$HOME/.anacron/spool"
0	*	*	*	*		XDG_RUNTIME_DIR=/run/user/$(id -u) notify-send -u critical -i "$HOME/.icon.png" "Friendly reminder" "Take a break and relax"
