# .crontab
0   *   *   *   *   /usr/sbin/anacron -t "$HOME/.anacron/anacrontab" -S "$HOME/.anacron/spool"
0   *   *   *   *   XDG_RUNTIME_DIR=/run/user/$(id -u) notify-send -u critical -i "$HOME/.icon.png" "Friendly reminder" "Take a break and relax"
0   10  *   *   5   XDG_RUNTIME_DIR=/run/user/$(id -u) notify-send -u critical -i "$HOME/.icon.png" "Lunch!" "Remember to order lunch"
